# services/xavier.service
[Unit]
Description=Xavier (Raspberry Pi client launcher)
Documentation=https://github.com/harmonicalabs/xavier
After=network-online.target otelcol.service
Wants=network-online.target
Requires=otelcol.service
StartLimitIntervalSec=0
StartLimitBurst=0

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/forge
ExecStart=/home/pi/forge/launch.sh

# Environment file for runtime configuration
EnvironmentFile=/etc/default/xavier

# Auto-restart on crash/failure
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xavier

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Resource limits
MemoryMax=512M
CPUQuota=100%

# Ensure service keeps running
KillMode=control-group
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
